# timers/01-minutely.target
[Unit]
Description=Minutely Timer Target
StopWhenUnneeded=yes
