# timers/01-usec.target
[Unit]
Description=Usec Timer Target
StopWhenUnneeded=yes
