# timers/02-hourly.timer
[Unit]
Description=Hourly Timer

[Timer]
OnBootSec=10min
OnCalendar=hourly
Unit=02-hourly.target
Persistent=true

[Install]
WantedBy=basic.target
