# timers/02-msec.target
[Unit]
Description=Msec Timer Target
StopWhenUnneeded=yes
