# timers/03-seconds.target
[Unit]
Description=Seconds Timer Target
StopWhenUnneeded=yes
