# timers/04-minutes.target
[Unit]
Description=Minutes Timer Target
StopWhenUnneeded=yes
