# timers/05-hours.timer
[Unit]
Description=Hours Timer

[Timer]
OnBootSec=10min
OnCalendar=hourly
Unit=05-hours.target
Persistent=true

[Install]
WantedBy=basic.target
