# timers/05-hours.target
[Unit]
Description=Hours Timer Target
StopWhenUnneeded=yes
