# timers/06-days.target
[Unit]
Description=Days Timer Target
StopWhenUnneeded=yes
