# timers/06-days.timer
[Unit]
Description=Daily Timer

[Timer]
OnBootSec=15min
OnCalendar=daily
Unit=06-days.target
Persistent=true

[Install]
WantedBy=basic.target
