# timers/06-quarterly.target
# /etc/systemd/system/analysis-06-quarterly.target
[Unit]
Description=Quarterly Timer Target
StopWhenUnneeded=yes
