# timers/07-semiannually.target
[Unit]
Description=Semi-annually Timer Target
StopWhenUnneeded=yes
