# timers/07-weeks.timer
[Unit]
Description=Weeks Timer

[Timer]
OnBootSec=20min
OnCalendar=weekly
Unit=07-weeks.target
Persistent=true

[Install]
WantedBy=basic.target
