# timers/07-weeks.target
[Unit]
Description=Weeks Timer Target
StopWhenUnneeded=yes
