# timers/08-months.target
[Unit]
Description=Months Timer Target
StopWhenUnneeded=yes
