# timers/08-months.timer
[Unit]
Description=Months Timer

[Timer]
OnBootSec=25min
OnCalendar=monthly
Unit=08-months.target
Persistent=true

[Install]
WantedBy=basic.target
