# timers/09-quarters.target
[Unit]
Description=Quarters Timer Target
StopWhenUnneeded=yes
