# timers/10-sixmonths.target
[Unit]
Description=Six Months Timer Target
StopWhenUnneeded=yes
