# timers/11-years.timer
[Unit]
Description=Years Timer

[Timer]
OnBootSec=40min
OnCalendar=yearly
Unit=11-years.target
Persistent=true

[Install]
WantedBy=basic.target
